# services/discord-bot.service
[Unit]
Description=BeatByte Discord Music Bot
After=network.target

[Service]
Type=simple
User=opc
WorkingDirectory=/home/opc/discord-bots/bots/music-bot
ExecStart=/usr/bin/node src/index.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
